# resources/cbox-cmdserver.service
[Unit]
Description=CBox Command Server
After=cbox-guestinit.service

[Service]
Type=simple
ExecStart=/usr/local/bin/cbox-cmdserver
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
